# homectrl.conf
# Ubuntu upstart file at /etc/init/homectrl.conf

start on runlevel [2345]
stop on runlevel [!2345]

script
	chdir /homectrl
	exec java -jar homectrl*.jar 
end script